# sv/zaqar/run
#!/bin/sh

sv start mongodb || exit 1

exec 2>&1
exec zaqar-server
